# ha_intercom/rootfs/etc/services.d/intercom/run
#!/command/with-contenv bash
set -euo pipefail
cd /opt/app
# Load env safely
set -a
source ./.env || true
set +a
exec npm start
