# ha_intercom/rootfs/etc/services.d/intercom/run
#!/command/with-contenv bashio
set -euo pipefail
cd /opt/app
set -a
[ -f .env ] && source .env || true
set +a
exec npm start
